# vendors/Moxa/da662a/romfsdisk/etc/rc.d/rc.local
# !/bin/sh
# Add you want to run daemon
